# C05_NodeJs_MongoDB/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --dbpath=/data/db
autorestart=false

[program:nodeapp]
command=node --loader esm index.mjs
directory=/app